# docker/etc/supervisor/conf.d/consul.conf
[program:consul]
command=consul agent -config-dir /etc/consul.d -join 192.168.150.10
stopsignal=INT
priority=1
