# docker/etc/supervisor/conf.d/rabbitmq.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
user=rabbitmq
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
