# docker/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
redirect_stderr=true
priority=1
stopsignal=TERM
